# learning/containerization/example2/test_data/alpine-3.19.1/var/run~Stashed changes/systemd/generator/local-fs.target.requires/boot-efi.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-601A\x2d21FB.service
After=systemd-fsck@dev-disk-by\x2duuid-601A\x2d21FB.service
After=blockdev@dev-disk-by\x2duuid-601A\x2d21FB.target

[Mount]
Where=/boot/efi
What=/dev/disk/by-uuid/601A-21FB
Type=vfat
