# learning/containerization/example2/test_data/alpine-3.19.1/var/run~Stashed changes/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
After=blockdev@dev-disk-by\x2did-dm\x2duuid\x2dLVM\x2dJPbMZB9GD3h6m8iGK0IjGjt2tWIkeWtsan0YwgTS3lEAXp85CFCSH8nUHSTLv9Yy.target

[Mount]
Where=/
What=/dev/disk/by-id/dm-uuid-LVM-JPbMZB9GD3h6m8iGK0IjGjt2tWIkeWtsan0YwgTS3lEAXp85CFCSH8nUHSTLv9Yy
Type=ext4
