# learning/containerization/example2/test_data/alpine-3.19.1/var/run/systemd/generator/boot.mount
# Automatically generated by systemd-fstab-generator

[Unit]
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
SourcePath=/etc/fstab
Before=local-fs.target
Requires=systemd-fsck@dev-disk-by\x2duuid-236af943\x2dae3a\x2d4cc0\x2d9985\x2d39ba3a1e2cb8.service
After=systemd-fsck@dev-disk-by\x2duuid-236af943\x2dae3a\x2d4cc0\x2d9985\x2d39ba3a1e2cb8.service
After=blockdev@dev-disk-by\x2duuid-236af943\x2dae3a\x2d4cc0\x2d9985\x2d39ba3a1e2cb8.target

[Mount]
Where=/boot
What=/dev/disk/by-uuid/236af943-ae3a-4cc0-9985-39ba3a1e2cb8
Type=ext4
