# learning/containerization/example2/test_data/alpine-3.19.1/var/run/user/1001/systemd/generator.late/xdg-desktop-autostart.target.wants/app-org.gnome.SettingsDaemon.DiskUtilityNotify@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/org.gnome.SettingsDaemon.DiskUtilityNotify.desktop
PartOf=graphical-session.target

Description=gnome-disk-utility notification plugin for GNOME Settings Daemon
After=graphical-session.target

[Service]
Type=exec
ExecStart=:/usr/libexec/gsd-disk-utility-notify
Restart=no
TimeoutSec=5s
Slice=app.slice
ExecCondition=/lib/systemd/systemd-xdg-autostart-condition "GNOME" ""
