# learning/containerization/example2/test_data/alpine-3.19.1/var/run/user/1001/systemd/generator.late/app-slack@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/home/roman/.config/autostart/slack.desktop
PartOf=graphical-session.target

Description=slack
After=graphical-session.target

[Service]
Type=exec
ExecStart=:/usr/bin/slack
Restart=no
TimeoutSec=5s
Slice=app.slice
