# learning/containerization/example2/test_data/alpine-3.19.1/var/run/user/1001/systemd/generator.late/app-org.gnome.DejaDup.Monitor@autostart.service
# Automatically generated by systemd-xdg-autostart-generator

[Unit]
Documentation=man:systemd-xdg-autostart-generator(8)
SourcePath=/etc/xdg/autostart/org.gnome.DejaDup.Monitor.desktop
PartOf=graphical-session.target

Description=Backup Monitor
After=graphical-session.target

[Service]
Type=exec
ExecStart=:/usr/libexec/deja-dup/deja-dup-monitor
Restart=no
TimeoutSec=5s
Slice=app.slice
